# airootfs/etc/skel/.config/autostart/transtheme.desktop
[Desktop Entry]
Name=Trans Theme
Exec=transtheme
Icon=utilities-terminal
Type=Application
Categories=System;Utitity